# etc/local.d/00hdapsd.start
#!/bin/sh

echo -1 > /sys/block/sdb/device/unload_heads
