# etc/local.d/10android.start
#!/bin/sh
/opt/android-sdk-update-manager/platform-tools/adb start-server
